# changes/crontab
* * * * * /home/changes/gen-pages >/dev/null
